# supervisord.conf
[supervisord]
nodaemon = true

[program:zmon-agent]
command = /run.sh
